# docker/postgis-tracts/postgis.sv.conf
[program:postgis]
user=postgres
command=/usr/pgsql-9.3/bin/postmaster -D  /var/lib/pgsql/9.3/data
stopsignal=INT
redirect_stderr=true
